# deploy/systemd/argus-health-check.timer
[Unit]
Description=ARGUS Health Check Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=30s
Persistent=true

[Install]
WantedBy=timers.target
