# deploy/systemd/argus-health-check.service
[Unit]
Description=ARGUS WAN Health Check Service
After=network.target argus-api.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/argus-cli --api-url http://localhost:8443 stats
User=argus
StandardOutput=journal
